# opencode.service
[Unit]
Description=OpenCode Server
After=network.target

[Service]
Type=simple
User=valle
Group=valle
WorkingDirectory=/home/valle/proyectos/opencode-bot

Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="OPENCODE_PORT=4096"

ExecStart=/usr/bin/opencode serve --port 4096 --hostname 0.0.0.0
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode

[Install]
WantedBy=multi-user.target
